# lottery-predictor/lottery-update.timer
[Unit]
Description=双色球数据每天自动更新定时器
Requires=lottery-update.service

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
Unit=lottery-update.service

[Install]
WantedBy=timers.target
